# cups/rootfs/etc/s6-overlay/s6-rc.d/dbus-daemon/run
#!/usr/bin/with-contenv bashio

mkdir -p /var/run/dbus
dbus-daemon --system --nofork
