# openwisp/rootfs/etc/services.d/openwisp-api/run
#!/usr/bin/with-contenv bashio
set -euo pipefail

# shellcheck disable=SC1091
. /run/openwisp/env

cd /opt/openwisp

export MODULE_NAME="api"
export DJANGO_SETTINGS_MODULE="openwisp.settings_api"
export CONTAINER_PORT="${API_APP_PORT}"
export UWSGI_PIDFILE="/opt/openwisp/uwsgi-api.pid"

envsubst < /opt/openwisp/uwsgi.conf.ini > /opt/openwisp/uwsgi-api.ini

exec uwsgi --ini /opt/openwisp/uwsgi-api.ini
