# openwisp/rootfs/etc/services.d/openwisp-celerybeat/run
#!/usr/bin/with-contenv bashio
set -euo pipefail

# shellcheck disable=SC1091
. /run/openwisp/env

cd /opt/openwisp

export DJANGO_SETTINGS_MODULE="openwisp.settings_dashboard"

exec celery -A openwisp beat -l "${DJANGO_LOG_LEVEL}" --schedule /data/celerybeat-schedule
