# openwisp/rootfs/etc/services.d/openwisp-celery-monitoring/run
#!/usr/bin/with-contenv bashio
set -euo pipefail

# shellcheck disable=SC1091
. /run/openwisp/env

if [ "${USE_OPENWISP_MONITORING}" != "True" ] || [ "${USE_OPENWISP_CELERY_MONITORING}" != "True" ]; then
  exec tail -f /dev/null
fi

cd /opt/openwisp

export DJANGO_SETTINGS_MODULE="openwisp.settings_dashboard"

exec celery -A openwisp worker -l "${DJANGO_LOG_LEVEL}" --queues monitoring \
  -n monitoring@%h --logfile /opt/openwisp/logs/celery_monitoring.log \
  ${OPENWISP_CELERY_MONITORING_COMMAND_FLAGS:-}
