# openwisp/rootfs/etc/services.d/openwisp-celery/run
#!/usr/bin/with-contenv bashio
set -euo pipefail

# shellcheck disable=SC1091
. /run/openwisp/env

cd /opt/openwisp

export DJANGO_SETTINGS_MODULE="openwisp.settings_dashboard"

exec celery -A openwisp worker -l "${DJANGO_LOG_LEVEL}" --queues celery \
  -n celery@%h --logfile /opt/openwisp/logs/celery.log \
  ${OPENWISP_CELERY_COMMAND_FLAGS:-}
