# openwisp/rootfs/etc/services.d/openwisp-websocket/run
#!/usr/bin/with-contenv bashio
set -euo pipefail

# shellcheck disable=SC1091
. /run/openwisp/env

cd /opt/openwisp

export MODULE_NAME="websocket"
export DJANGO_SETTINGS_MODULE="openwisp.settings_websocket"
export DJANGO_WEBSOCKET_HOST="${DJANGO_WEBSOCKET_HOST:-0.0.0.0}"
export CONTAINER_PORT="${WEBSOCKET_APP_PORT}"

exec daphne -b "${DJANGO_WEBSOCKET_HOST}" -p "${WEBSOCKET_APP_PORT}" --proxy-headers openwisp.asgi:application --access-log -
