# cups/rootfs/etc/s6-overlay/s6-rc.d/cups-server/run
#!/usr/bin/with-contenv bashio

until [ -e /var/run/avahi-daemon/socket ]; do
  sleep 1s
done

cupsd -f

