# openwisp/rootfs/etc/supervisor/conf.d/dashboard.conf
[program:dashboard]
command=gunicorn --bind 0.0.0.0:8000 --workers 2 --timeout 120 openwisp.wsgi
directory=/opt/openwisp
user=openwisp
environment=MODULE_NAME=dashboard
autostart=true
autorestart=true
priority=30
depends_on=postgres,redis
stdout_logfile=/var/log/openwisp-dashboard.log
stderr_logfile=/var/log/openwisp-dashboard-error.log