# openwisp/rootfs/etc/supervisor/conf.d/redis.conf
[program:redis]
command=/usr/bin/redis-server /etc/redis/redis.conf
user=redis
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/redis.log
stderr_logfile=/var/log/redis-error.log