# openwisp/rootfs/etc/services.d/openwisp-dashboard/run
#!/usr/bin/with-contenv bashio
set -euo pipefail

# shellcheck disable=SC1091
. /run/openwisp/env

cd /opt/openwisp

export MODULE_NAME="dashboard"
export DJANGO_SETTINGS_MODULE="openwisp.settings_dashboard"
export CONTAINER_PORT="${DASHBOARD_APP_PORT}"
export UWSGI_PIDFILE="/opt/openwisp/uwsgi-dashboard.pid"

envsubst < /opt/openwisp/uwsgi.conf.ini > /opt/openwisp/uwsgi-dashboard.ini

exec uwsgi --ini /opt/openwisp/uwsgi-dashboard.ini
