# openwisp/rootfs/etc/services.d/openwisp-celery-firmware/run
#!/usr/bin/with-contenv bashio
set -euo pipefail

# shellcheck disable=SC1091
. /run/openwisp/env

if [ "${USE_OPENWISP_FIRMWARE}" != "True" ] || [ "${USE_OPENWISP_CELERY_FIRMWARE}" != "True" ]; then
  exec tail -f /dev/null
fi

cd /opt/openwisp

export DJANGO_SETTINGS_MODULE="openwisp.settings_dashboard"

exec celery -A openwisp worker -l "${DJANGO_LOG_LEVEL}" --queues firmware_upgrader \
  -n firmware_upgrader@%h --logfile /opt/openwisp/logs/celery_firmware_upgrader.log \
  ${OPENWISP_CELERY_FIRMWARE_COMMAND_FLAGS:-}
